# service/AmazonBot/run
#!/bin/sh
BOT_USER=ytnobody
cd /home/ytnobody/AmazonBot
exec 2>&1
exec setuidgid $BOT_USER envdir ./env ./worker.pl
sleep 5 ### reconnect interval
